# systemd/inform-me.service
[Unit]
Description=Service for informing me

[Service]
Type=dbus
BusName=local.inform.GTalk
ExecStart=/usr/bin/python3 /usr/bin/info-service.py
